# base/config/kafka/supervisor/zookeeper.ini
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
